# server-setup/obsidian-git-sync.service
[Unit]
Description=obsidian-git-sync daemon (hermes vault watcher)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/obsidian-git-sync/daemon
# node 절대경로 필수 — systemd 는 PATH 를 보장하지 않는다. `which node` 결과로 치환할 것.
ExecStart=/usr/bin/node /home/ubuntu/obsidian-git-sync/daemon/dist/index.js
# hermes 는 bare repo 와 같은 박스이므로 REMOTE 는 SSH 가 아닌 로컬 경로.
Environment=VAULT_PATH=/home/ubuntu/wiki
Environment=REMOTE=/home/ubuntu/vault.git
Environment=DEVICE_ID=hermes
Environment=DEBOUNCE_MS=3000
Environment=AUTOSAVE_IDLE_MS=300000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
